# files/etc/uci-defaults/99-tollgate-setup
#!/bin/sh
# TollGate first-boot configuration

# 1. Firewall configuration
if ! uci -q get firewall.tollgate_rules; then
    uci add firewall include
    uci set firewall.@include[-1].path='/etc/config/firewall-tollgate'
    uci rename firewall.@include[-1]='tollgate_rules'
    uci commit firewall
fi

# 2. Configure uhttpd listener ports
current_listen=$(uci -q get uhttpd.main.listen_http)
if [ "$current_listen" != "0.0.0.0:8080" ]; then
    uci set uhttpd.main.listen_http='0.0.0.0:8080'
    uci add_list uhttpd.main.listen_http='[::]:8080'
    uci commit uhttpd
fi

# 3. Configure DNS settings
current_dns=$(uci -q get network.lan.dns)
if [ "$current_dns" = "127.0.0.1" ] || [ "$current_dns" = "::1" ] || [ -z "$current_dns" ]; then
    # Only change DNS if it's set to localhost or not set
    uci -q delete network.lan.dns
    uci add_list network.lan.dns='1.1.1.1'  # CloudFlare primary DNS
    uci add_list network.lan.dns='1.0.0.1'  # CloudFlare secondary DNS
fi
uci set network.lan.domain='lan'

# Configure dnsmasq
uci set dhcp.@dnsmasq[0].domainneeded='1'
uci set dhcp.@dnsmasq[0].localise_queries='1' 
uci set dhcp.@dnsmasq[0].rebind_protection='1'
uci set dhcp.@dnsmasq[0].local='/lan/'
uci set dhcp.@dnsmasq[0].domain='lan'
uci set dhcp.@dnsmasq[0].expandhosts='1'
uci commit network
uci commit dhcp

# 4. Configure NoDogSplash
if ! uci -q get nodogsplash.@nodogsplash[0]; then
    # Create nodogsplash section if it doesn't exist
    uci add nodogsplash nodogsplash
fi

# Set basic NoDogSplash configuration
uci set nodogsplash.@nodogsplash[0].enabled='1'
uci set nodogsplash.@nodogsplash[0].gatewayname='TollGate Nodogsplash'
uci set nodogsplash.@nodogsplash[0].gatewayinterface='br-lan'

# Ensure TollGate protocol port is allowed
uci add_list nodogsplash.@nodogsplash[0].users_to_router='allow tcp port 2121'
uci add_list nodogsplash.@nodogsplash[0].users_to_router='allow tcp port 8080'
uci commit nodogsplash

# 5. Add TollGate port to firewall (if not in firewall-tollgate)
if ! grep -q "port 2121" /etc/config/firewall-tollgate 2>/dev/null; then
    # Add rule directly to main firewall if not in tollgate-specific file
    uci add firewall rule
    uci set firewall.@rule[-1].name='Allow-TollGate-Protocol'
    uci set firewall.@rule[-1].src='wan'
    uci set firewall.@rule[-1].proto='tcp'
    uci set firewall.@rule[-1].dest_port='2121'
    uci set firewall.@rule[-1].target='ACCEPT'
    uci commit firewall
fi

# 6. Add first-login-setup hook to profile if it doesn't exist
if ! grep -q "first-login-setup" /etc/profile; then
    # Append the first-login-setup code to the end of profile
    cat >> /etc/profile << 'EOF'

# TollGate first login setup
if [ ! -f /etc/first_login_done ] && [ -t 0 ] && [ -t 1 ]; then
    /usr/local/bin/first-login-setup
fi
EOF
    echo "Added first-login-setup hook to /etc/profile"
fi

# Ensure script exits successfully
exit 0